# systemd/psntui-sync.timer
[Unit]
Description=psnTUI - Periodic PSN trophy sync (every 4 hours)
Documentation=https://github.com/TripShuti/psntui

[Timer]
OnCalendar=*-*-* 00/3:50:00
Persistent=true
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
