# systemd/psntui-sync.timer
[Unit]
Description=psnTUI - Periodic PSN trophy sync (every 4 hours)
Documentation=https://github.com/TripShuti/psntui

[Timer]
OnBootSec=15min
OnUnitActiveSec=4h
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
